# infra/deploy/systemd/rfsite-cloud-api.service
[Unit]
Description=RF Site Telemetry Cloud API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/rfsite/cloud-api

# Ajustá path al venv real
ExecStart=/opt/rfsite/.venv/bin/rfsite-cloud-api --config /etc/rfsite-cloud-api/config.yaml

Restart=always
RestartSec=2

# hardening razonable
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rfsite /tmp /var/log/rfsite

[Install]
WantedBy=multi-user.target
